# scripts/launchd/com.kkk4oru.finish-position-predict.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent for the daily finish-position-prediction docker pipeline.

  WHY MAC LAUNCHD (NOT CLOUDFLARE CRON):
    Cloudflare Containers reap batch instances at ~90-110 s regardless of
    sleepAfter, but the DuckDB feature build + per-category CatBoost / XGBoost
    scoring needs ~10 min. The Cloudflare Cron Trigger has therefore been
    disabled in apps/finish-position-cron/wrangler.jsonc (triggers.crons = []),
    and this LaunchAgent drives the proven local docker pipeline instead —
    same Dockerfile (apps/finish-position-predict-container/Dockerfile), same
    image (finish-position-predict-local:split2), same env vars
    (SOURCE_DATABASE_URL = local Colima PG, NEON_DATABASE_URL = Neon write
    target, RUN_DATE = today JST in YYYYMMDD form), idempotent UPSERT into
    race_finish_position_model_predictions.

  SCHEDULE (Phase-3 Fix #1 — JRA serve-skew fix):
    Two StartCalendarInterval entries fire this plist:

    1. JST 03:00 daily  — NAR + Ban-ei only (PREDICT_CATEGORIES=nar,ban-ei).
       jvd_se (JRA mirror) is not yet available at 03:00 JST (mirror arrives
       ~09:03 JST), so attempting JRA at this time always finds races=0 and
       wastes ~30s. Scoping to nar+ban-ei eliminates the wasted work and makes
       the 03:00 run category-clean.

    2. JST 09:30 daily  — ALL categories (PREDICT_CATEGORIES unset = all).
       By 09:30 JST the PC-Keiba daily sync has mirrored jvd_se (~09:03 JST)
       and the hot worker has 12+ h of JRA advance-odds snapshots in D1
       (advance-odds window opens at 19:00 JST prior evening). This run is the
       first opportunity for JRA predictions with REAL odds instead of the
       OOD-median fallback that was applied when only the 03:00 run existed.
       On weekdays (no JRA races) the JRA category finds races=0 and exits in
       ~30s — harmless. On Sat/Sun it finds the day's JRA races and scores them
       with real odds.

    The race-prediction-guard (10:00-20:40 JST, 20-min cadence) continues to
    re-kick the pipeline during race hours for bataiju / odds freshness.

  CATCH-UP BEHAVIOUR:
    StartCalendarInterval will fire on next wake if the scheduled time was
    missed while the Mac was asleep (this is launchd's "catch up on missed
    runs" behaviour). The wrapper script's UPSERT is idempotent, so a delayed
    fire just overwrites today's predictions safely.

  STDOUT/STDERR PATHS:
    launchd does NOT expand "~", so the paths below are literal absolute
    paths under /Users/kkk4oru/Library/Logs/finish-position-predict/. The
    wrapper script also writes a date-stamped log to the same dir.

  SECRETS:
    No secrets live in this plist. The wrapper script reads
    NEON_DATABASE_URL from apps/local-postgresql/.env.replica at run time and
    masks it in the dated log.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kkk4oru.finish-position-predict</string>

  <key>Program</key>
  <string>/Users/kkk4oru/ghq/github.com/kkkaoru/horse-racing-data/scripts/launchd/finish-position-predict-daily.sh</string>

  <key>WorkingDirectory</key>
  <string>/Users/kkk4oru/ghq/github.com/kkkaoru/horse-racing-data</string>

  <key>StartCalendarInterval</key>
  <array>
    <!-- JST 03:00 daily — NAR + Ban-ei only (JRA mirror not yet available). -->
    <dict>
      <key>Hour</key>
      <integer>3</integer>
      <key>Minute</key>
      <integer>0</integer>
    </dict>
    <!-- JST 09:30 daily — ALL categories. JRA mirror arrives ~09:03 JST;
         hot-worker D1 has 12+ h of JRA advance-odds snapshots by this time.
         On weekdays JRA finds races=0 and skips in ~30s (harmless).
         On Sat/Sun this is the first run that can fetch REAL JRA odds. -->
    <dict>
      <key>Hour</key>
      <integer>9</integer>
      <key>Minute</key>
      <integer>30</integer>
    </dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/kkk4oru/Library/Logs/finish-position-predict/stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kkk4oru/Library/Logs/finish-position-predict/stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- Only non-secret defaults belong here. The wrapper script reads any
         credentials it needs (NEON_DATABASE_URL, SOURCE_DATABASE_URL) from
         apps/local-postgresql/.env.replica + a hardcoded localhost URL. -->
    <key>RUN_DATE_MODE</key>
    <string>auto</string>
    <key>PREDICT_DAYS_AHEAD</key>
    <string>0</string>
    <!-- PATH so the wrapper finds docker / colima / awk / date without
         relying on the user's interactive shell init. /opt/homebrew/bin is
         Apple-Silicon Homebrew (docker, colima, bun all live there). -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
