# scripts/launchd/com.kkk4oru.race-prediction-guard.plist
<?xml version="1.0" encoding="UTF-8"?>
<!--
  LaunchAgent for the race-prediction freshness guard.

  WHY MAC LAUNCHD (NOT CLOUDFLARE CRON):
    Cloudflare Containers reap batch instances at ~90-110 s regardless of
    sleepAfter (see com.kkk4oru.finish-position-predict.plist for the full
    backstory). The finish-position docker pipeline this guard re-invokes is
    the same ~4 min workload, so the same launchd-on-Mac escape applies.
    The running-style worker is itself Cloudflare-Worker friendly, but the
    coordination ("kick only if today is incomplete / fresh") needs to run
    somewhere that can talk to both Neon and D1 — that is most cleanly here.

  SCHEDULE (~46 fires/day in JST):
    * Hour 0..9   (Minute 0)        -> guard TODAY    JST (10 fires, hourly)
    * Hour 10..20 (Minute 0,20,40)  -> guard TODAY    JST (33 fires, 20-min cadence)
    * Hour 21..23 (Minute 0)        -> guard TODAY + TOMORROW JST (3 fires, hourly)
    macOS launchd uses the device's local timezone, and this machine is
    assumed to already run JST. If you ever change the system timezone you
    MUST adjust each Hour key accordingly — launchd has no timezone-aware
    option. The wrapper script computes the target JST date independently
    via `date -u -v+9H`, so the *which-day-to-guard* logic is robust to TZ
    drift; only the *when-to-fire* depends on the system TZ.

  WHY 20-MIN CADENCE DURING JST 10-20 (RACE HOURS):
    Bataiju (馬体重) is announced ~T-30..40 min before post and saved to D1.
    The finish-position pipeline takes ~4 min measured (max 5 min conserv.).
    With hourly cadence, bataiju could land 1 min after a trigger, leaving
    the next trigger only 1 min before post — not enough time to complete.
    With 20-min cadence, worst case: bataiju lands 1 min after a trigger;
    the next trigger fires 19 min later (e.g. at T-11 for T-30 bataiju),
    pipeline completes at T-6. All 12 Oi races on 2026-06-10 (14:00-20:50)
    have ≥15 min margin under this cadence (verified: tightest is R1 with
    15 min, best is R2/R6 with 30 min).
    The 0..9 morning band stays hourly (no race-day post times; the 03:00
    JST primary cron handles overnight predictions). The 21..23 pre-warm
    stays hourly (no race-day posts; only overnight data refreshes).

  FRESHNESS RE-PREDICTION (guard script change):
    During JST 10-20 the guard's finish-position skip-when-complete logic
    is bypassed: even when fp_actual >= expected_count the pipeline is
    re-kicked so predictions incorporate the latest bataiju/odds. The
    concurrent-run lock (/tmp/finish-position-predict.lock) is still
    respected — two docker runs never race each other.

  CATCH-UP / IDEMPOTENCY:
    StartCalendarInterval queues missed firings while the Mac is asleep and
    fires once on next wake. The guard re-evaluates Neon coverage every tick;
    Neon is the state of truth, so a duplicate fire after a kick simply sees
    "complete" and exits. No state file is needed.

  STDOUT/STDERR PATHS:
    launchd does NOT expand "~", so the paths below are literal absolute
    paths under /Users/kkk4oru/Library/Logs/race-prediction-guard/. The
    wrapper script also writes a per-JST-day dated log to the same dir.

  SECRETS:
    No secrets live in this plist. The wrapper script reads
    REALTIME_ADMIN_TOKEN from apps/sync-realtime-data/.dev.vars and
    NEON_DATABASE_URL from apps/local-postgresql/.env.replica at run time,
    and masks both in the dated log.
-->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.kkk4oru.race-prediction-guard</string>

  <key>Program</key>
  <string>/Users/kkk4oru/ghq/github.com/kkkaoru/horse-racing-data/scripts/launchd/race-prediction-guard.sh</string>

  <key>WorkingDirectory</key>
  <string>/Users/kkk4oru/ghq/github.com/kkkaoru/horse-racing-data</string>

  <key>StartCalendarInterval</key>
  <array>
    <!-- 00:00..09:00 JST: guard TODAY JST (10 fires, hourly). -->
    <dict><key>Hour</key><integer>0</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>1</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>2</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>3</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>4</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>5</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>6</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>7</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>8</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>9</integer><key>Minute</key><integer>0</integer></dict>
    <!-- 10:00..20:40 JST: guard TODAY JST (33 fires, 20-min race-hours cadence).
         Bataiju lands ~T-30..40 min; 20-min cadence guarantees a trigger fires
         within the [bataiju, post-10min] window for every race. Pipeline run
         time ~4 min (max 5 min conserv.), so completion is always ≥5 min before
         post. Worst-case margin on today's Oi card (14:00-20:50): 15 min (R1). -->
    <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Hour</key><integer>10</integer><key>Minute</key><integer>40</integer></dict>
    <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Hour</key><integer>11</integer><key>Minute</key><integer>40</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Hour</key><integer>12</integer><key>Minute</key><integer>40</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Hour</key><integer>13</integer><key>Minute</key><integer>40</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Hour</key><integer>14</integer><key>Minute</key><integer>40</integer></dict>
    <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Hour</key><integer>15</integer><key>Minute</key><integer>40</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Hour</key><integer>16</integer><key>Minute</key><integer>40</integer></dict>
    <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Hour</key><integer>17</integer><key>Minute</key><integer>40</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Hour</key><integer>18</integer><key>Minute</key><integer>40</integer></dict>
    <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Hour</key><integer>19</integer><key>Minute</key><integer>40</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>20</integer></dict>
    <dict><key>Hour</key><integer>20</integer><key>Minute</key><integer>40</integer></dict>
    <!-- 21:00..23:00 JST: guard TODAY + TOMORROW JST (3 fires, hourly pre-warm). -->
    <dict><key>Hour</key><integer>21</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>22</integer><key>Minute</key><integer>0</integer></dict>
    <dict><key>Hour</key><integer>23</integer><key>Minute</key><integer>0</integer></dict>
  </array>

  <key>RunAtLoad</key>
  <false/>

  <key>KeepAlive</key>
  <false/>

  <key>StandardOutPath</key>
  <string>/Users/kkk4oru/Library/Logs/race-prediction-guard/stdout.log</string>

  <key>StandardErrorPath</key>
  <string>/Users/kkk4oru/Library/Logs/race-prediction-guard/stderr.log</string>

  <key>EnvironmentVariables</key>
  <dict>
    <!-- PATH so the wrapper finds bunx / wrangler / jq / docker / colima / uv
         / date without relying on the user's interactive shell init.
         /opt/homebrew/bin is Apple-Silicon Homebrew (where bun, jq, uv,
         docker, colima all live). -->
    <key>PATH</key>
    <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
</dict>
</plist>
