# apps/pc-keiba-viewer/scripts/com.kkkaoru.win5-overlay.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.kkkaoru.win5-overlay</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/zsh</string>
        <string>-lc</string>
        <string>cd /Users/kkk4oru/ghq/github.com/kkkaoru/horse-racing-data/apps/pc-keiba-viewer &amp;&amp; /opt/homebrew/bin/bun run src/scripts/generate-win5-overlay.ts</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>DATABASE_URL_LOCAL</key>
        <string>postgresql://horse_racing:horse_racing@127.0.0.1:15432/horse_racing</string>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Run at 09:00 JST every Saturday and Sunday (WIN5 race days). -->
    <key>StartCalendarInterval</key>
    <array>
        <dict>
            <key>Weekday</key>
            <integer>6</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
        <dict>
            <key>Weekday</key>
            <integer>0</integer>
            <key>Hour</key>
            <integer>9</integer>
            <key>Minute</key>
            <integer>0</integer>
        </dict>
    </array>

    <key>StandardOutPath</key>
    <string>/tmp/win5-overlay-stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/win5-overlay-stderr.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
